# wzz/systemd/user/flatpak-session-helper.service
[Unit]
Description=flatpak session helper

[Service]
BusName=org.freedesktop.Flatpak
ExecStart=/usr/libexec/flatpak-session-helper
Type=dbus
