# wzz/systemd/user/flatpak-portal.service
[Unit]
Description=flatpak portal

[Service]
BusName=org.freedesktop.portal.Flatpak
ExecStart=/usr/libexec/flatpak-portal
Type=dbus
