# wzz/systemd/user/vino-server.service
[Unit]
Description=Vino VNC server

[Service]
Type=dbus
BusName=org.gnome.Vino
ExecStart=/usr/libexec/vino-server
Restart=on-abnormal
